# udev/99-battery-monitor.rules
SUBSYSTEM=="power_supply", ATTR{online}=="*", ATTR{type}=="Mains", RUN+="/opt/battery-monitor/battery-event-handler.sh"
