# udev/99-battery-monitor.rules
SUBSYSTEM=="power_supply", ATTR{online}=="*", ATTR{type}=="Mains", RUN+="/usr/local/bin/battery-event-handler.sh"
